# systemd/check-public-ip.service
[Unit]
Description=Chequeo de IP Pública

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/local/bin/check-public-ip.sh
